# kubelet/1.27/kubelet_1.27.9/DEBIAN/postinst
#!/bin/sh -e

if [ "$1" = "configure" ] ; then
        # Initial installation
        systemctl preset kubelet.service >/dev/null 2>&1 || :
fi
